# OpenComputer/opencomputer/service/templates/com.opencomputer.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{label}</string>
    <key>ProgramArguments</key>
    <array>
        <string>{executable}</string>
        <string>--headless</string>
        <string>--profile</string>
        <string>{profile}</string>
        <string>gateway</string>
    </array>
    <key>WorkingDirectory</key>
    <string>{workdir}</string>
    <key>RunAtLoad</key>
    <true/>
    <!--
      KeepAlive=dict (not <true/>): only respawn on crash (exit != 0).
      With KeepAlive=<true/>, launchd respawns even after intentional
      `kill -9`, making `oc gateway stop` impossible without first
      `launchctl bootout`. The Hermes plist uses this same dict form;
      we adopted it after the 2026-05-08 respawn-loop incident
      (drainer SQLite error every 1s for 33h while user couldn't kill
      the daemon). See gateway/server.py docstring for full context.
    -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    <key>StandardOutPath</key>
    <string>{stdout_log}</string>
    <key>StandardErrorPath</key>
    <string>{stderr_log}</string>
</dict>
</plist>
