# OpenComputer/opencomputer/service/templates/opencomputer-profile-analyze.timer
[Unit]
Description=OpenComputer daily profile-analyze
Requires=opencomputer-profile-analyze.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
